# rc.local
/usr/sbin/tun.sh || exit 1                                                          
exit 0
